# bash/systemd/podman-backup-weekly.timer
[Unit]
Description=Run Podman weekly backup on Sunday at 3 AM
Requires=podman-backup-weekly.service

[Timer]
# Run every Sunday at 3:00 AM
OnCalendar=Sun *-*-* 03:00:00
Persistent=true

# Add randomized delay to avoid exact concurrent runs
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
